# strixhalo/systemd/user/1bit-watchdog.service
[Unit]
Description=1bit-watchdog — upstream tracker (24h soak, auto-merge to bong fork)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.cargo/bin/1bit-watchdog check
Environment=RUST_LOG=info
Nice=10
