# strixhalo/systemd/user/1bit-tier-mint.service
[Unit]
Description=1bit-tier-mint — BTCPay + Patreon webhook → Premium tier JWT issuance
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Secrets live in a root-owned env file outside the repo. This unit file
# is tracked; the env file is NOT. See docs/wiki/tier-jwt-flow.md §Secrets.
EnvironmentFile=%h/.config/1bit/tier-mint.env
ExecStart=%h/.cargo/bin/1bit-tier-mint
Restart=on-failure
RestartSec=5s
Nice=10

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native

[Install]
WantedBy=default.target
