# crates/1bit-power/packaging/1bit-halo-power-log.timer
# /etc/systemd/system/1bit-halo-power-log.timer
#
# Fire 1bit-halo-power-log.service every 30 s. AccuracySec kept generous so
# systemd can batch wakeups with neighbours.
[Unit]
Description=1bit-halo-power — metrics sampler (30 s cadence)
Documentation=https://1bit.systems/1bit-halo-power

[Timer]
OnBootSec=30s
OnUnitActiveSec=30s
AccuracySec=2s
Unit=1bit-halo-power-log.service
Persistent=false

[Install]
WantedBy=timers.target
