# crates/1bit-power/packaging/halo-power-log.service
# /etc/systemd/system/halo-power-log.service
#
# Emit one JSON metrics line to the journal. Triggered by halo-power-log.timer.
[Unit]
Description=halo-power — emit one metrics JSON line to journald
Documentation=https://1bit.systems/halo-power

[Service]
Type=oneshot
ExecStart=/usr/local/bin/halo-power log
# Metric sampling reads /sys/class/hwmon only — no capability bump needed.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
# Tag the log lines so `journalctl -t halo-power` filters just us.
SyslogIdentifier=halo-power
