# strixhalo/systemd/1bit-halo-archive.timer
[Unit]
Description=Nightly 1bit-halo-archive run (strixhalo → Pi additive sync)

[Timer]
OnCalendar=*-*-* 03:00:00
# If strixhalo was powered down at 3am, run on next boot once we're online.
Persistent=true
# Randomize by up to 20 min to avoid racing other cron-shaped work.
RandomizedDelaySec=20min
Unit=1bit-halo-archive.service

[Install]
WantedBy=timers.target
