# strixhalo/systemd/1bit-halo-gh-trio.timer
[Unit]
Description=1bit systems GitHub automation — quartermaster, magistrate, librarian

[Timer]
OnBootSec=3min
OnUnitActiveSec=30min
Persistent=true
Unit=1bit-halo-gh-trio.service

[Install]
WantedBy=timers.target
