# strixhalo/systemd/1bit-halo-beta-expire.timer
[Unit]
Description=Daily 03:00 timer for 1bit-halo-beta-expire (10-day beta TTL)

[Timer]
# Nightly at 03:00 local. Pairs well with 1bit-halo-archive at 03:00:00 +20min
# jitter — the expire sweep is the earlier cheap job.
OnCalendar=*-*-* 03:00:00
# Don't hammer on boot; wait 5 min for services + mesh to settle before
# catching up on a missed run.
OnBootSec=5min
# If strixhalo was off during a 03:00 slot, run the missed sweep once
# we're back up.
Persistent=true
RandomizedDelaySec=5min
Unit=1bit-halo-beta-expire.service

[Install]
WantedBy=timers.target
