# strixhalo/systemd/1bit-halo-kokoro.service
[Unit]
Description=1bit systems kokoro-server (TTS) — echo_mouth backend
After=network.target
Wants=network.target

[Service]
Type=simple
WorkingDirectory=__HOME__/repos/1bit-halo-kokoro
ExecStart=/usr/bin/bun run __HOME__/repos/1bit-halo-kokoro/server.ts
Restart=on-failure
RestartSec=5
Environment=HALO_KOKORO_BIN=__HOME__/repos/1bit-halo-kokoro/kokoro_tts
Environment=HALO_KOKORO_RESOURCES=__HOME__/repos/1bit-halo-kokoro/resources
Environment=HALO_KOKORO_PORT=8083
Environment=HALO_KOKORO_HOST=127.0.0.1

[Install]
WantedBy=default.target
