# strixhalo/systemd/1bit-halo-bitnet.service
[Unit]
Description=1bit systems bitnet decode server
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/bitnet_decode "__HOME__/1bit systems/models/halo-1bit-2b.h1b" --server 8080
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
