# strixhalo/systemd/1bit-halo-anvil.service
[Unit]
Description=1bit-halo-anvil — rebuild rocm-cpp on new commit, bench, post to Discord
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=__HOME__/bin/1bit-anvil.sh
Nice=10
IOSchedulingClass=idle
StandardOutput=journal
StandardError=journal
