# strixhalo/systemd/1bit-halo-patreon.service
[Unit]
Description=1bit-watch-patreon — one-shot Patreon member poller, Quartermaster dispatch, Discord announce
Documentation=https://github.com/bong-water-water-bong/1bit-systems/blob/main/crates/1bit-agents/bin/1bit-watch-patreon.rs
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Installed via `cargo install --path crates/1bit-agents --bin 1bit-watch-patreon`.
# Binary is read-only against the Patreon API (GET /members only).
ExecStart=/home/bcloud/.cargo/bin/1bit-watch-patreon
# PATREON_ACCESS_TOKEN + ECHO_BOT_TOKEN + HALO_DISCORD_ANNOUNCE_CHANNEL
# live in a drop-in file so the tracked unit stays secret-free.
# Create /home/bcloud/.config/systemd/user/1bit-halo-patreon.service.d/env.conf
# with:
#   [Service]
#   Environment=PATREON_ACCESS_TOKEN=...
#   Environment=ECHO_BOT_TOKEN=...
#   Environment=HALO_DISCORD_ANNOUNCE_CHANNEL=<numeric channel id>
# Without the token drop-in the binary logs "no token, skipping" and exits 0.
Environment=RUST_LOG=onebit_watch_patreon=info
Nice=10
IOSchedulingClass=idle
StandardOutput=journal
StandardError=journal
