# strixhalo/systemd/1bit-halo-memory-sync.timer
[Unit]
Description=Run 1bit-halo-memory-sync every 15 minutes
Requires=1bit-halo-memory-sync.service

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
