# strixhalo/systemd/1bit-tts.service
[Unit]
Description=1bit-tts — Qwen3-TTS C++/ggml server (OpenAI /v1/audio/speech)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=RUST_LOG=info
ExecStart=%h/.local/bin/1bit-tts-server \
    -m %h/.local/share/1bit/models/qwen3-tts-0p6b-base/model.gguf \
    -v %h/.local/share/1bit/models/qwen3-tts-tokenizer/model.gguf \
    -H 127.0.0.1 -p 8095
Restart=on-failure
RestartSec=5s
Nice=5

[Install]
WantedBy=default.target
