# strixhalo/systemd/1bit-halo-whisper-vulkan.service
[Unit]
# Vulkan-backed whisper-server variant. Intended for sliger (Intel Arc B580,
# Xe2) where whisper.cpp is built with `-DGGML_VULKAN=ON` into
# /opt/whisper-vulkan/. Binds :8190 (the reserved audio-pipeline port) and
# serves large-v3 weights — large dedicated GPU lets us skip the base-model
# fallback the strixhalo unit uses.
#
# On the audio host (sliger), install this unit + drop the RDNA3 / Strix
# Halo `1bit-halo-whisper.service` variant. Clients on strixhalo call
# sliger.local:8190 over the headscale mesh.
Description=1bit-halo-whisper (Vulkan) — STT on Arc B580 / Xe2 (:8190)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=WHISPER_PREFIX=/opt/whisper-vulkan
Environment=WHISPER_MODEL=/opt/whisper-vulkan/models/ggml-large-v3.bin
ExecStart=/opt/whisper-vulkan/bin/whisper-server \
    --host 0.0.0.0 \
    --port 8190 \
    --model ${WHISPER_MODEL} \
    --threads 8 \
    --no-timestamps
Restart=on-failure
RestartSec=5
Nice=0

# Audio pipeline wants low jitter — inherit the real-time scheduler class
# if available (CachyOS ships with kernel.sched_rt_runtime_us tuned).
CPUSchedulingPolicy=rr
CPUSchedulingPriority=50

# Xe2 GPU access. `render` group for DRM, `video` for Vulkan context.
SupplementaryGroups=render video

[Install]
WantedBy=default.target
