# strixhalo/systemd/halo-gh-trio.service
[Unit]
Description=halo-ai GitHub automation runner — walks issues, PRs, commits

[Service]
Type=oneshot
ExecStart=__HOME__/bin/halo-quartermaster.sh
ExecStart=__HOME__/bin/halo-magistrate.sh
ExecStart=__HOME__/bin/halo-librarian.sh
Nice=10
IOSchedulingClass=idle
StandardOutput=journal
StandardError=journal
