# strixhalo/systemd/halo-archive.service
[Unit]
Description=halo-archive — additive rsync of halo-ai + repos from strixhalo to Pi ZFSPool
# Mesh must be up; if the Pi is offline, the unit fails gracefully and the
# timer retries at the next schedule (not a reboot event).
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=oneshot
Nice=10
IOSchedulingClass=idle
ExecStart=__HOME__/bin/halo-archive.sh
# Never succeed silently on mesh failure — systemd logs the non-zero exit
# and the next timer fire still tries. No retry-now.
SuccessExitStatus=0
StandardOutput=journal
StandardError=journal
