# strixhalo/systemd/halo-whisper.service
[Unit]
Description=halo-ai whisper-server (STT) — echo_ear backend
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=__HOME__/repos/whisper.cpp/build/bin/whisper-server \
    --host 127.0.0.1 \
    --port 8082 \
    --model __HOME__/repos/whisper.cpp/models/ggml-base.en.bin \
    --threads 8
Restart=on-failure
RestartSec=5
Nice=10

[Install]
WantedBy=default.target
