# strixhalo/systemd/halo-gh-trio.service
[Unit]
Description=1bit systems GitHub automation runner — walks issues, PRs, commits

[Service]
Type=oneshot
ExecStart=__HOME__/bin/1bit-quartermaster.sh
ExecStart=__HOME__/bin/1bit-magistrate.sh
ExecStart=__HOME__/bin/1bit-librarian.sh
Nice=10
IOSchedulingClass=idle
StandardOutput=journal
StandardError=journal
