# strixhalo/systemd/halo-sd.service
[Unit]
Description=halo-ai stable-diffusion (native HIP) server
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=__HOME__/repos/stable-diffusion.cpp/build-hip/bin/sd-server \
    --model __HOME__/models/sdxl/sd_xl_base_1.0.safetensors \
    --listen-ip 127.0.0.1 \
    --listen-port 8081
Restart=on-failure
RestartSec=5
Environment=HSA_OVERRIDE_GFX_VERSION=11.5.1
Nice=5

[Install]
WantedBy=default.target
