# strixhalo/systemd/halo-agent.service
[Unit]
Description=1bit systems agent runtime
After=halo-bitnet.service
Wants=halo-bitnet.service

[Service]
Type=simple
ExecStart=/usr/local/bin/agent_cpp
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
