# strixhalo/systemd/halo-lemonade.service
[Unit]
Description=halo-ai lemonade-server (lemond daemon) on :8000
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/home/bcloud/.local/bin/lemond --host 127.0.0.1 --port 8000 /home/bcloud/.cache/lemonade
Restart=on-failure
RestartSec=5
Nice=10

[Install]
WantedBy=default.target
