# strixhalo/systemd/halo-beta-expire.service
[Unit]
Description=halo-beta-expire — 10-day TTL enforcer for bearers + beta memory files
# We don't strictly need the network here (revoke talks to the local
# Headscale + local Caddy), but if mesh is down we'd rather fail fast
# than half-apply.
After=network-online.target caddy.service
Wants=network-online.target

[Service]
Type=oneshot
Nice=15
IOSchedulingClass=idle
# Apply mode: the timer is the scheduled enforcer, not a preview.
ExecStart=%h/repos/halo-workspace/strixhalo/bin/halo-beta-expire.sh --apply
# A partial failure (one handle fails to revoke) exits 1. We log that and
# let the timer retry next fire; no retry-now storm.
SuccessExitStatus=0
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
