# strixhalo/systemd/halo-gh-trio.timer
[Unit]
Description=halo-ai GitHub automation — quartermaster, magistrate, librarian

[Timer]
OnBootSec=3min
OnUnitActiveSec=30min
Persistent=true
Unit=halo-gh-trio.service

[Install]
WantedBy=timers.target
