# strixhalo/systemd/strix-watch-github.timer
[Unit]
Description=halo-watch-github poller (read-only issue + PR classifier)

[Timer]
# 5 min after boot, then every 5 min. Matches the binary's default
# HALO_GH_POLL_SECONDS=300 and its lookback window of 2×poll_seconds.
OnBootSec=5min
OnUnitActiveSec=5min
Persistent=true
Unit=strix-watch-github.service

[Install]
WantedBy=timers.target
