# strixhalo/systemd/strix-lemonade.service
[Unit]
Description=strix-ai-rs 1bit-lemonade (OpenAI-compat model gateway on :8200)
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=__HOME__/.cargo/bin/1bit-lemonade --bind 127.0.0.1:8200
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=onebit_lemonade=info
Nice=15

[Install]
WantedBy=default.target
