# strixhalo/systemd/strix-cloudflared.service
[Unit]
Description=Cloudflare Tunnel for api.halo-ai.studio (halo-server :8180)
After=network-online.target strix-server.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/cloudflared tunnel --config %h/.cloudflared/config.yml run
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
