# strixhalo/systemd/strix-landing.service
[Unit]
Description=1bit-landing (marketing page + live status on :8190)
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=%h/.local/bin/1bit-landing
Restart=on-failure
RestartSec=5
Environment=ONEBIT_LOG=info
Nice=15

[Install]
WantedBy=default.target
