# strixhalo/systemd/strix-echo.service
[Unit]
Description=1bit-echo WebSocket voice server (:8085)
After=network-online.target strix-server.service strix-kokoro.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/1bit-echo --bind 127.0.0.1:8085 --codec opus
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
