# strixhalo/systemd/strix-landing.service
[Unit]
Description=strix-ai-rs halo-landing (marketing page + live status on :8190)
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=__HOME__/.cargo/bin/halo-landing
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=halo_landing=info
Nice=15

[Install]
WantedBy=default.target
