# strixhalo/systemd/strix-server.service
[Unit]
Description=strix-ai-rs 1bit-server (gen-2 Rust, real HIP backend on gfx1151)
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=__HOME__/.local/bin/1bit-server-real --bind 127.0.0.1:8180 --model "__HOME__/1bit systems/models/halo-1bit-2b.h1b"
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=onebit_server=info,onebit_router=info
Environment=LD_LIBRARY_PATH=__HOME__/repos/rocm-cpp/build:/opt/rocm/lib
Nice=10

[Install]
WantedBy=default.target
