# strixhalo/systemd/strix-watch-github.service
[Unit]
Description=halo-watch-github — read-only GitHub poller for halo-agents specialists
Documentation=https://github.com/strix-ai-rs/halo-workspace/blob/main/docs/wiki/Crate-halo-agents-watch-github.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Installed via `cargo install --path crates/halo-agents --bin halo-watch-github`.
# Binary is read-only: no GitHub writes (no comments / labels / reactions).
ExecStart=/home/bcloud/.cargo/bin/halo-watch-github
# GITHUB_TOKEN lives in a drop-in file so the tracked unit stays secret-free.
# Create /home/bcloud/.config/systemd/user/strix-watch-github.service.d/token.conf
# with:
#   [Service]
#   Environment=GITHUB_TOKEN=ghp_...
# Without the drop-in the binary logs "no token, skipping" and exits 0.
Environment=RUST_LOG=info
Nice=10
IOSchedulingClass=idle
StandardOutput=journal
StandardError=journal
