# strixhalo/systemd/strix-watch-discord.service
[Unit]
Description=halo-watch-discord — Discord presence for 1bit-agents specialists (lurker, mention-reply only)
After=network-online.target
Wants=network-online.target
# Not strictly required — the bot degrades gracefully if 1bit-server is
# absent (status replies just show "?"). Left as a soft hint.
Wants=strix-server.service strix-landing.service

[Service]
Type=simple
# Installed via `cargo install --path crates/1bit-agents` or the
# `halo install 1bit-agents` path (packages.toml).
ExecStart=%h/.cargo/bin/halo-watch-discord
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=RUST_LOG=halo_watch_discord=info
# DISCORD_BOT_TOKEN is NOT defined here — drop it into a root-protected
# override at:
#   ~/.config/systemd/user/strix-watch-discord.service.d/token.conf
# with:
#   [Service]
#   Environment="DISCORD_BOT_TOKEN=…"
# HALO_DISCORD_CHANNELS belongs in the same override file.
Nice=15
# Defensive: if the token is missing the binary exits 0, which would
# otherwise send systemd into a tight "simple" restart loop at Restart=
# on-failure semantics. on-failure ignores clean exits, so we're fine —
# the unit goes inactive and the operator is expected to provide a token.

[Install]
# Deliberately NOT enabled by default. Operator must run:
#   systemctl --user enable --now strix-watch-discord
# after dropping in the token override. Fresh-box installs should never
# bring this unit up on their own.
WantedBy=default.target
